# server/crontab
* * * * * /usr/local/bin/ddns-fwu.py /root/etc/dfwu.ini # DDNS Firewall Update